# home/nand/.config/systemd/user/ipfs.service
[Unit]
Description=IPFS daemon
After=network.target

[Service]
ExecStart=/usr/bin/ipfs daemon

[Install]
WantedBy=multiuser.target
