# home/nand/.config/systemd/user/notmuch.service
[Unit]
Description=Update notmuch database

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new
TimeoutStartSec=30s
