# home/nand/.config/systemd/user/notmuch.timer
[Unit]
Description=Automatically update notmuch database

[Timer]
OnActiveSec=1s
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
