# etc/systemd/system/munin-node.service.d/override.conf
[Unit]
After=openvpn-client@tina.service

[Service]
Restart=on-failure
RestartSec=10s
