# etc/cron.monthly/zfs-scrub
#!/bin/sh
POOL=yggdrasil
/sbin/zpool scrub ${POOL}
