# home/nand/.config/systemd/user/pulseaudio.service.d/override.conf
[Service]
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50
Nice=-10
