# home/nand/.config/systemd/user/pulseaudio.service.d/override.conf
[Service]
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=50
CPUAffinity=7,23
Nice=-10
Environment=LADSPA_PATH=/usr/lib64/ladspa
