# etc/systemd/system/mpd.service.d/override.conf
[Service]
ProtectKernelModules=no
ExecStart=/usr/bin/ionice -c realtime /usr/bin/mpd --no-daemon
