# etc/cron.weekly/adblock2privoxy
#!/bin/sh
sudo -u privoxy -g privoxy /usr/bin/adblock2privoxy -t /etc/privoxy/ab2p.task
systemctl restart privoxy
