# etc/cron.weekly/le-renew-webroot.sh
#!/bin/sh
/usr/local/sbin/le-renew-webroot >> /var/log/le-renewal.log
