# etc/cron.weekly/zfs-scrub
#!/bin/sh
POOL=zarathustra
/sbin/zpool scrub ${POOL}
