# etc/systemd/system/force-cache.service
[Unit]
Description=Force stuff into ARC cache

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/force-cache
