# etc/systemd/system/zfs-force-cache.service
[Unit]
Description=Force stuff into ARC cache

[Service]
Type=oneshot
ExecStart=/usr/bin/ionice -c idle /usr/local/sbin/zfs-force-cache
