# etc/systemd/system/bitcoind.service.d/override.conf
[Unit]
Wants=tor.service
After=tor.service
